# java-collections-ru/advanced-streams/src/test/resources/fixtures/s1.conf
[program:options]
environment="X_FORWARDED_variable=value,  "

key=value
